# configs/automatic-backup/squid-proxy-deployment-backup.service
[Unit]
Description=Squid Deployment Backup

[Service]
ExecStart=/usr/bin/squid-proxy-deployment-backup.sh
